# home-manager/scripts/scripts/disable-screensaver.service
[Unit]
After=mysql.service

[Service]
ExecStart=/usr/local/bin/disable-screensaver.sh

[Install]
WantedBy=default.target
